# deploy/segrada.service
[Unit]
Description=Workout Ledger app
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/segrada/app
ExecStart=/opt/segrada/venv/bin/gunicorn --workers 2 --bind 127.0.0.1:8020 wsgi:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
